# picpickr.conf
[program:myserver]
command=gunicorn run:app -w 1 -b 0.0.0.0:80

[supervisord]
logfile=picpickr.log
loglevel=debug
user=root